# ops/ansible/roles/node/templates/node.service.j2
# {{ ansible_managed | default('Autogenerated by Ansible') }}
# Systemd unit for Animica node managed via Docker Compose.
# Templated variables consumed here (with sane defaults):
# - animica_compose_dir (required)
# - node_compose_file (default: node.docker-compose.yml)
# - docker_compose_cmd  (default: /usr/bin/docker compose)
# - node_service_name   (default: animica-node)
# - node_service_description (default)
# - animica_user / animica_group (default: animica)
# - node_rpc_port_http  (default: 8545)
# - node_limit_nofile   (default: 262144)

[Unit]
Description={{ node_service_description | default('Animica Node (Docker Compose)') }}
Documentation=https://github.com/animica-labs
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
ConditionPathExists={{ animica_compose_dir }}/{{ node_compose_file | default('node.docker-compose.yml') }}

[Service]
Type=oneshot
RemainAfterExit=yes
User={{ animica_user | default('animica') }}
Group={{ animica_group | default('animica') }}
WorkingDirectory={{ animica_compose_dir }}
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Increase file descriptor limits for P2P workloads
LimitNOFILE={{ node_limit_nofile | default(262144) }}
# Validate and pre-pull images (best-effort)
ExecStartPre=/usr/bin/bash -lc "{{ docker_compose_cmd | default('/usr/bin/docker compose') }} -f {{ animica_compose_dir }}/{{ node_compose_file | default('node.docker-compose.yml') }} config -q"
ExecStartPre=/usr/bin/bash -lc "{{ docker_compose_cmd | default('/usr/bin/docker compose') }} -f {{ animica_compose_dir }}/{{ node_compose_file | default('node.docker-compose.yml') }} pull --quiet || true"
# Start the stack
ExecStart={{ docker_compose_cmd | default('/usr/bin/docker compose') }} -f {{ animica_compose_dir }}/{{ node_compose_file | default('node.docker-compose.yml') }} up -d --remove-orphans
# Optional health wait (HTTP /healthz on local port)
ExecStartPost=/usr/bin/bash -lc "for i in {1..60}; do curl -sf http://127.0.0.1:{{ node_rpc_port_http | default(8545) }}/healthz >/dev/null && exit 0; sleep 2; done; exit 1"
# Graceful stop keeps volumes/data
ExecStop={{ docker_compose_cmd | default('/usr/bin/docker compose') }} -f {{ animica_compose_dir }}/{{ node_compose_file | default('node.docker-compose.yml') }} stop
# Reload will re-apply the compose (pull optional)
ExecReload=/usr/bin/bash -lc "{{ docker_compose_cmd | default('/usr/bin/docker compose') }} -f {{ animica_compose_dir }}/{{ node_compose_file | default('node.docker-compose.yml') }} pull --quiet || true"
ExecReload={{ docker_compose_cmd | default('/usr/bin/docker compose') }} -f {{ animica_compose_dir }}/{{ node_compose_file | default('node.docker-compose.yml') }} up -d
# If the compose command itself fails, retry a bit
Restart=on-failure
RestartSec=5s
SyslogIdentifier={{ node_service_name | default('animica-node') }}

[Install]
WantedBy=multi-user.target
